# deploy/systemd/quicast-wttest-cert-sync.timer
[Unit]
Description=Check for renewed wttest.quicast.de TLS certificate

[Timer]
OnBootSec=10m
OnUnitActiveSec=6h
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
